# images/modoboa/rootfs/lib/supervisor/conf/dcron.conf
[program:dcron]
# Base configuration
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Command to run
command=crond -f -s /etc/cron.d -M /usr/bin/dcron-mailer
