# images/modoboa/rootfs/lib/supervisor/conf/modoboa.conf
[program:modoboa]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=modoboa
group=modoboa
directory=/data/instance
command=gunicorn -c /lib/modoboa/gunicorn.conf.py instance.wsgi:application
