# images/modoboa/rootfs/lib/supervisor/conf/syslog.conf
[program:syslog]
# Base configuration
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Command to run
command=syslogd -n -r --no-unixaf --no-klog
