# images/postfix/rootfs/etc/runit/services/postfix/run
#!/bin/sh

exec /usr/libexec/postfix/master -d
